# switch_service.service
[Unit]
Description=Switches between the fallback and current OS partition

[Service]
Type=simple
ExecStart=/home/root/switch.sh

[Install]
WantedBy=multi-user.target
