# switch_service.service
[Unit]
Description=Switches between the fallback and current OS partition
RequiresMountsFor=/home/root /var/lib/uboot

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/bin/bash /home/root/switch.sh
TimeoutSec=infinity

[Install]
WantedBy=multi-user.target